# lab-trading-dashboard/scripts/api-signals.service
# systemd unit for api_signals.py (Python signals API)
# Serves /api/calculate-signals, /api/open-position, /api/sync-open-positions.
# Node (lab-trading-dashboard) proxies to this; set PYTHON_SIGNALS_URL=http://127.0.0.1:5001 for Node.
#
# Install on cloud (one time):
#   sudo cp api-signals.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable api-signals
#   sudo systemctl start api-signals
#
# Ensure Node has in /etc/lab-trading-dashboard.secrets.env:
#   PYTHON_SIGNALS_URL=http://127.0.0.1:5001
# Then: sudo systemctl restart lab-trading-dashboard

[Unit]
Description=LAB Trading Dashboard Python Signals API (api_signals.py)
After=network.target postgresql.service
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/lab-trading-dashboard/python
ExecStart=/usr/bin/python3 api_signals.py
Restart=always
RestartSec=10
Environment=PYTHON_SIGNALS_PORT=5001
EnvironmentFile=-/etc/lab-trading-dashboard.secrets.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=api-signals

[Install]
WantedBy=multi-user.target
